# tools/traffic_tool/crontab
5 * * * * php /root/resetSSR/traffic_statistics.php >> /root/resetSSR/ss.log
